# overlay/etc/cron.d/ushahidi
*/5  * * * * www-data /usr/bin/php /var/www/ushahidi/platform/artisan datasource:outgoing >/dev/null
*/5  * * * * www-data /usr/bin/php /var/www/ushahidi/platform/artisan datasource:incoming >/dev/null
*/5  * * * * www-data /usr/bin/php /var/www/ushahidi/platform/artisan savedsearch:sync >/dev/null
*/5  * * * * www-data /usr/bin/php /var/www/ushahidi/platform/artisan notification:queue >/dev/null
*/5  * * * * www-data /usr/bin/php /var/www/ushahidi/platform/artisan webhook:send >/dev/null
